# ops/twin-sync/bot-vps/systemd/twin-runtime-bridge.timer
[Unit]
Description=Run Twin Runtime Bridge every 5 minutes

[Timer]
OnBootSec=45s
OnUnitActiveSec=5min
AccuracySec=10s
Unit=twin-runtime-bridge.service
Persistent=true

[Install]
WantedBy=timers.target
