# ops/scout-worker/scout-vps/systemd/scout-poll.path
[Unit]
Description=Watch Scout inbox and trigger poll immediately

[Path]
PathExistsGlob=/srv/scout/inbox/requests/*.json
Unit=scout-poll.service

[Install]
WantedBy=multi-user.target
