# ops/scout-worker/scout-vps/systemd/scout-poll.timer
[Unit]
Description=Run Scout poll every minute

[Timer]
OnCalendar=*:0/1
Persistent=true

[Install]
WantedBy=timers.target
