# ops/separate-sanitizer/sanitizer-vps/systemd/sanitizer-sanitize.path
[Unit]
Description=Watch sanitizer raw inbox and trigger sanitize immediately

[Path]
PathExistsGlob=/srv/sanitizer/inbox/raw/*.json
Unit=sanitizer-sanitize.service

[Install]
WantedBy=multi-user.target
