# ops/separate-sanitizer/sanitizer-vps/systemd/sanitizer-sanitize.service
[Unit]
Description=Sanitizer raw -> checked pipeline
After=network.target

[Service]
Type=oneshot
User=sanitizer
Group=sanitizer
ExecStart=/usr/local/bin/sanitizer_service.py
WorkingDirectory=/srv/sanitizer

[Install]
WantedBy=multi-user.target
