# ops/separate-sanitizer/bot-vps/systemd/sanitizer-import.path
[Unit]
Description=Watch signed bundle inbox and trigger import immediately

[Path]
PathModified=/home/claudeclaw/import/inbox
Unit=sanitizer-import.service

[Install]
WantedBy=multi-user.target
