# ops/scout-worker/bot-vps/systemd/adjutant-scout-request-push.timer
[Unit]
Description=Run Adjutant->Scout request push every minute

[Timer]
OnCalendar=*:0/1
Persistent=true

[Install]
WantedBy=timers.target
